# deployment/twit.conf
[program:twit]
command=/home/ubuntu/twit/venv/bin/gunicorn -b localhost:8000 -w 4 twit:app
directory=/home/ubuntu/twit
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
